# scripts/crons/example.crontab
02 * * * * socorro /data/socorro/scripts/crons/cron_aggregates.sh
05 07 * * * socorro /data/socorro/scripts/crons/cron_builds.sh
11 * * * * socorro /data/socorro/scripts/crons/cron_bugzilla.sh
20 02 * * 0 socorro /data/socorro/scripts/crons/cron_clean_tcbs.sh
05 21 * * 2 socorro /data/socorro/scripts/crons/cron_create_partitions.sh
15 00 * * * socorro /data/socorro/scripts/crons/cron_daily_crashes.sh
*/5 * * * * socorro /data/socorro/scripts/crons/cron_status.sh
00 02 * * * socorro /data/socorro/scripts/crons/cron_tcurl.sh
