# config/package/cron.d/socorro
*/5 * * * * socorro /data/socorro/application/scripts/crons/crontabber.sh
